# templates/server/ubuntu/lib/systemd/system/xxx.service
# XXX systemd service file

[Unit]
StartLimitBurst=5
StartLimitIntervalSec=30

[Service]
TimeoutSec=infinity
Restart=on-failure
RestartSec=5
